# deploy/pkm-backup.timer
[Unit]
Description=Daily encrypted PKM backup timer

[Timer]
OnCalendar=*-*-* 03:15:00 UTC
RandomizedDelaySec=20m
Persistent=true
Unit=pkm-backup.service

[Install]
WantedBy=timers.target
